# bluealsa-aplay.service
# /etc/systemd/system/bluealsa-aplay.service
#
# Joue l'audio Bluetooth vers le loopback ALSA (substream 1).
# bluealsa-aplay se reconnecte automatiquement à chaque appareil A2DP.

[Unit]
Description=BlueALSA aplay — Bluetooth audio vers hw:Loopback,0,1
After=bluealsa.service bluetooth-init.service sound.target
Requires=bluealsa.service

[Service]
Type=simple
# bt_sink is defined in /etc/asound.conf:
# normalizes BT format (44100/48000 Hz) → 48000 Hz S16LE before writing to the loopback
ExecStart=/usr/bin/bluealsa-aplay --pcm=bt_sink
Restart=always
RestartSec=5
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target
