# carplay-audio.service
[Unit]
Description=carplay-audio DSP daemon
Documentation=https://github.com/e46-bt/daemon
After=network.target bluetooth.target shairport-sync.service bluealsa.service bluealsa-aplay.service
Wants=shairport-sync.service bluealsa.service bluealsa-aplay.service

[Service]
Type=simple
ExecStart=/usr/local/bin/carplay-audio
Restart=on-failure
RestartSec=2

# RuntimeDirectory creates /run/carplay-audio/ for the Unix socket
RuntimeDirectory=carplay-audio
RuntimeDirectoryMode=0755

# CAP_SYS_NICE allows SCHED_FIFO on audio threads without running as root
AmbientCapabilities=CAP_SYS_NICE
CapabilityBoundingSet=CAP_SYS_NICE
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
